# kiwi-descriptions/root/etc/systemd/system/github-actions-remote-executor.service
[Unit]
Description=GitHub Actions Remote Executor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/github-actions-remote-executor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-actions-remote-executor

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
